# cookbooks/rightscale/files/centos/yum-makecache.cron
#!/bin/bash

# Updates the metadata cache for yum
yum makecache -y
